# chub/files/var/spool/cron/crontabs/root
# Cron tasks for SliTaz Community HUB aka chub
#

# Update RRD graphs images
*/5  * * * * /usr/bin/chub rrd >/dev/null
