# chub/files/var/spool/cron/crontabs/root
# Cron tasks for SliTaz Community HUB aka chub
#

# Update RRD images
*/5 * * * * /home/slitaz/www/chub/lib/makegraphs >/dev/null
